# setup/scrapper/scrapper.service
[Unit]
Description=TIPE-BSSR Scrapper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tipe
WorkingDirectory=/opt/tipe-bssr
ExecStart=/opt/tipe-bssr/venv/bin/python -m src.main scrapper --data-dir /opt/tipe-bssr/data
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
